# scripts/deploy/alpha-trader-api.service
[Unit]
Description=Alpha Trader API (Fastify)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/alpha-trader
EnvironmentFile=-/opt/alpha-trader/.env
Environment=HOST=0.0.0.0
Environment=PORT=3000
Environment=NODE_ENV=production
Environment=SERVE_WEB_APP=true
# Node 22+ required (ESM workspace libs). Prefer nvm default, then system node.
Environment=NVM_DIR=/root/.nvm
ExecStart=/bin/bash -lc 'source "${NVM_DIR}/nvm.sh" 2>/dev/null && nvm use 22 >/dev/null 2>&1; exec node /opt/alpha-trader/apps/alpha-trader-server/dist/main.js'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target